# openresty-rtmp/supervisor_mp4_push.conf
[program:ffmpeg_push]
directory=/
command=/usr/local/bin/ffmpeg -re -stream_loop -1 -i /www/app/test.mp4 -c copy -f flv rtmp://localhost:1935/live/film
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log